# dot_config/systemd/user/hyprland-session.target
[Unit]
Description=Hyprland user session
BindsTo=hyprland-session.target
Wants=graphical-session-pre.target
After=graphical-session-pre.target
