# dot_config/systemd/user/hyprland-session2.target
[Unit]
Description=Hyprland user session
BindsTo=graphical-session.target
After=graphical-session-pre.target
Wants=graphical-session-pre.target
