# dot_config/systemd/user/seadrive.service
[Unit]
Description=Seadrive Gui Client
PartOf=hyprland-session.target
After=hyprland-session.target

[Service]
ExecStart=/usr/bin/seadrive-gui
Restart=on-failure

[Install]
WantedBy=hyprland-session.target
