# iso/archiso/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
# ExecStart=-/sbin/agetty --autologin liveuser --noclear %I 38400 linux
ExecStart=-/sbin/agetty --autologin root --noclear %I 38400 linux
